# deb_repo/thirdreality-openHAB/cmake/openhab/prerm
#!/bin/sh

set -e

case "$1" in
remove)
    systemctl stop openhab.service
    systemctl disable openhab.service
    systemctl revert openhab.service
    ;;
esac
