# deb_repo/thirdreality-openHAB/deb-package/cmake/openhab/postinst
#!/bin/sh

set -e

adduser --system --no-create-home --group --disabled-login openhab

chown -hR openhab:openhab /opt/openhab

systemctl enable openhab.service
