# deb_repo/thirdreality-openHAB/cmake/zigbee2mqtt/postinst
#!/bin/bash

set -e

case "$1" in
configure)
    systemctl enable zigbee2mqtt.service
    ;;
esac
