# deb_repo/thirdreality-openHAB/deb-package/cmake/zigbee2mqtt/prerm
#!/bin/sh

set -e

systemctl stop zigbee2mqtt.service
systemctl disable zigbee2mqtt.service
systemctl revert zigbee2mqtt.service
